# systemd/.config/systemd/user/ydotool.socket
[Unit]
Description=YDotool socket

[Socket]
ListenStream=/tmp/.ydotool_socket
SocketMode=0600

[Install]
WantedBy=sockets.target
